# deploy/systemd/checkal-reset-pausa-llm.service
# Limpa a flag PAUSA_LLM à meia-noite (o teto diário de custo LLM recomeça).
# NÃO toca em nenhum gate de compliance — só a flag de custo.
[Unit]
Description=CheckAL — reset da flag PAUSA_LLM (teto diário de custo)

[Service]
Type=oneshot
ExecStart=/usr/bin/rm -f /run/checkal/PAUSA_LLM
