# deploy/polaris/units/checkal-gestor.timer
[Unit]
Description=CheckAL agente gestor
[Timer]
OnCalendar=*-*-* 07:15
RandomizedDelaySec=180
Persistent=true
AccuracySec=60
Unit=checkal-agente@gestor.service
[Install]
WantedBy=timers.target
