# deploy/polaris/units/checkal-maestro-digest.timer
[Unit]
Description=CheckAL agente maestro-digest
[Timer]
OnCalendar=*-*-* 07:50
RandomizedDelaySec=180
Persistent=true
AccuracySec=60
Unit=checkal-agente@maestro-digest.service
[Install]
WantedBy=timers.target
